# meteoapi.service
[Unit]
Description=meteodata

[Service]
WorkingDirectory=/home/pi/meteodata/myapp/
ExecStartPre=/usr/bin/npm install
ExecStart=/usr/bin/npm run start
